# daemons/service/bumper/run
#!/bin/sh
exec /var/bumper/bumper-daemon.py
